# .config/cconf/environments/browser-start-page/browser-start-page.service
[Unit]
Description=Local browser start page

[Service]
ExecStart=/usr/bin/python3 -m http.server 8420 --bind 127.0.0.1 --directory %h/.local/share/browser-start-page
Restart=on-failure

[Install]
WantedBy=default.target
